# my_config/autostart/battery-watch.desktop
[Desktop Entry]
Type=Application
Name=Battery Watch
Exec=/home/meoww/.local/bin/battery-watch.sh
X-GNOME-Autostart-enabled=true

